# home/pi/.config/systemd/user/powerpoint-usb@.service
[Unit]
Description="Launch Powerpoint on USB plug"

[Service]
ExecStart=%h/.local/bin/powerpoint-usb.sh %i Management_visuel_methode.pps

